# capella-dockerimages/remote/supervisord.conf
[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
user=techuser
autorestart=true

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
user=techuser
autorestart=true
environment=DISPLAY=":10"

[supervisord]
nodaemon=true
childlogdir=/tmp